# systemd/air-scan-pull.service
[Unit]
Description=Air Scan Pull Scanner (OpenWrt router)
After=network.target air-scan-api.service
Wants=network.target

[Service]
Type=simple
User=matheau
WorkingDirectory=/home/matheau/code/air_scan
ExecStart=/home/matheau/miniforge3/bin/python3 scanners/pull_scanner.py --interval 60
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
