# systemd/air-scan-sync.service
[Unit]
Description=Air Scan Known Device Sync (port_scan integration)
After=network.target air-scan-api.service
Wants=network.target

[Service]
Type=simple
User=matheau
WorkingDirectory=/home/matheau/code/air_scan
ExecStart=/home/matheau/miniforge3/bin/python3 scanners/sync_known_devices.py --interval 300
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
