# scanners/99-mobile-sync.rules
ACTION=="change", SUBSYSTEM=="net", KERNEL=="eth*", ATTR{operstate}=="up", RUN+="/bin/systemctl --no-block start mobile-sync.service"
