# scanners/mobile-ble-scanner-test.service
[Unit]
Description=Mobile BLE scanner (test mode — no recording)
After=bluetooth.service hciuart.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/matheau/air_scan
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 /home/matheau/air_scan/scanners/mobile_ble_scanner.py --no-record

Restart=no

StandardOutput=journal
StandardError=journal

# No [Install] section — never enable at boot, always started manually
